# etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/usr/sbin/sshd
autostart=true
autorestart=true
redirect_stderr=true