# etc/supervisor/conf.d/sshd.conf
[program:sshd]
command=/etc/init.d/ssh
autostart=true
autorestart=true
redirect_stderr=true